# deploy/supervisor/poisapay-horizon.conf
; Horizon queue workers for PoisaPay (non-Docker / bare-metal or VM deploy).
; Install: copy to /etc/supervisor/conf.d/ then:
;   supervisorctl reread && supervisorctl update && supervisorctl start poisapay-horizon
[program:poisapay-horizon]
process_name=%(program_name)s
command=php /var/www/poisapay/artisan horizon
directory=/var/www/poisapay
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/poisapay/storage/logs/horizon.log
stopwaitsecs=3600
; Graceful shutdown so in-flight jobs finish before restart
stopsignal=TERM
