# deploy/supervisor/poisapay-scheduler.conf
; Laravel scheduler for PoisaPay (chain-tick, chain-health, accrue-credit).
; Alternative to a crontab entry; keeps schedule:work supervised.
[program:poisapay-scheduler]
process_name=%(program_name)s
command=php /var/www/poisapay/artisan schedule:work
directory=/var/www/poisapay
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/poisapay/storage/logs/scheduler.log
stopsignal=TERM
